# pyftpdlib.conf
description "Ftps server using python pyftpdlib"

start on runlevel [2345]

exec /home/ubuntu/ftps.sh